# oc_hack.service
[Unit]
Description=OC Hack Service

[Service]
ExecStart=/etc/oc_hack/service.sh

[Install]
WantedBy=multi-user.target
